# bases/overlay-openrc/etc/init.d/scw-set-hostname
#!/sbin/openrc-run
# Copyright (c) 2015-2018 Scaleway <opensource@scaleway.com>
# Released under the MIT license.

description="Set SCW hostname"

depend() {
    need root
    after net
}

start() {
    ebegin "Setting SCW hostname"
    /usr/local/sbin/scw-set-hostname
    echo "hostname=\"$(cat /etc/hostname)\"" > /etc/conf.d/hostname
    eend $?
}
