# bases/overlay-openrc/etc/init.d/scw-sync-kernel-extra
#!/sbin/runscript
# Copyright (c) 2015 Scaleway <opensource@scaleway.com>
# Released under the MIT license.

description="Downloads kernel modules/headers from SCW' mirror"

depend() {
    need root
    before kmod-static-nodes
}

start() {
    ebegin "Syncing kernel modules/headers"
    /usr/local/sbin/scw-sync-kernel-modules
    eend $?
}
