# bases/overlay-openrc/etc/init.d/scw-swapfile
#!/sbin/runscript
# Copyright (c) 2015 Scaleway <opensource@scaleway.com>
# Released under the MIT license.

description="Enables a swapfile"

depend() {
    need root
    before kmod-static-nodes
}

start() {
    ebegin "Starting swapfile"
    /usr/local/sbin/scw-swapfile start
    eend $?
}

stop() {
    ebegin "Stopping swapfile"
    /usr/local/sbin/scw-swapfile stop
    eend $?
}
