# skeleton-openrc/etc/init.d/disconnect-extra-volumes
#!/sbin/runscript
# Copyright (c) 2015 Scaleway <opensource@scaleway.com>
# Released under the MIT license.

description="Disconnects additional SCW volumes"

depend() {
    after mount-ro fsck
    before nbd-root-disconnect
}

start() {
    ebegin "Disconnecting extra volumes"
    /usr/local/sbin/oc-nbd-disconnect-extra
    eend $?
}
