# skeleton-openrc/etc/init.d/sync-connect-extra-volumes
#!/sbin/runscript
# Copyright (c) 2015 Scaleway <opensource@scaleway.com>
# Released under the MIT license.

description="Waits for additional SCW volumes"

depend() {
    need fsck
}

start() {
    ebegin "Connecting NBD extra volumes"
    /usr/local/sbin/oc-sync-connect-extra-volumes
    eend $?
}
