# bases/overlay-systemd/etc/systemd/system/scw-generate-net-config.service
[Unit]
Description=Generate network configuration
Before=networking.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/scw-generate-net-config

[Install]
WantedBy=networking.service
