# image-builder/patches/etc/rc.local

if [ ! -f ~/.ssh/id_rsa ]; then
    ssh-keygen  -f ~/.ssh/id_rsa -t rsa  -N  ''
fi

exit 0

